# supervisord.conf
[supervisord]
nodaemon=true

[program:openvpn]
command=/bin/bash -c "/usr/sbin/openvpn --config /etc/openvpn/openvpn.conf"