# supervisord.conf
[supervisord]
nodaemon=true

[program:openvpn]
command=service openvpn start